# sensehat/claudometer-sensehat.service
# systemd unit for the Claudometer Sense HAT display (headless, auto-start on boot).
#
# Install (adjust User and paths to match your Pi):
#   1. Clone the repo to /home/pi/claudometer
#   2. Build the venv run.sh uses (so apt's sense_hat stays visible):
#        cd /home/pi/claudometer && uv venv --system-site-packages && uv pip install requests
#   3. sudo cp sensehat/claudometer-sensehat.service /etc/systemd/system/
#   4. sudo systemctl daemon-reload
#   5. sudo systemctl enable --now claudometer-sensehat
#   6. journalctl -u claudometer-sensehat -f   # watch logs
#
# The service user must be in the 'input' group to read the Sense HAT joystick
# (sudo usermod -aG input pi) and typically 'i2c'/'gpio' for the LED matrix.

[Unit]
Description=Claudometer Sense HAT usage display
Wants=network-online.target
After=network-online.target multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/claudometer
ExecStart=/home/pi/claudometer/sensehat/run.sh
Restart=on-failure
RestartSec=10
# Credentials live here by default; override if you moved them.
# Environment=CLAUDE_MONITOR_DATA=/home/pi/claudometer/data
# Add display options as needed, e.g. rotation to match orientation:
# Environment=
# ExecStart=/home/pi/claudometer/sensehat/run.sh --rotation 180

[Install]
WantedBy=multi-user.target
